# scripts/systemd/insider-ibkr-paper-trader.service
[Unit]
Description=Insider Trades IBKR paper trader
After=insider-ibgateway.service
Requires=insider-ibgateway.service

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/insider_trader
EnvironmentFile=/etc/insider_trades.env
ExecStart=/home/opc/insider_trader/scripts/vm/run_ibkr_paper_trader.sh
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
